# LED_Control/LEDOff.service
[Unit]
Description="Turns off the 4 user LEDs"

[Service]
type=simple
ExecStart=/usr/bin/TurnOffLeds.sh

[Install]
WantedBy=multi-user.target
